# deploy/asena-main.service
# ─── ASENA MAIN SERVICE ───
# NEXUS OASIS AI Operating System - Ana Orkestrasyon
# 7/24 arka planda çalışan ana servis (Watcher + Orchestrator + Swarm)
#
# Kurulum:
#   sudo cp deploy/asena-main.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable asena-main
#   sudo systemctl start asena-main
#
# Yönetim:
#   sudo systemctl status asena-main
#   sudo journalctl -u asena-main -f
#   sudo systemctl restart asena-main

[Unit]
Description=ASENA Main - NEXUS OASIS AI Operating System (Orchestrator + Swarm)
Documentation=https://github.com/nexus-oasis/asena
After=network.target network-online.target asena-vgate.service asena-gateway.service
Wants=network-online.target
Requires=asena-vgate.service asena-gateway.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/ASENA_CORE

# Ortam değişkenleri
Environment="RUST_LOG=info"
Environment="ASENA_ENV=production"
Environment="ASENA_DATA_DIR=/root/ASENA_CORE/data"
Environment="ASENA_VGATE_URL=http://127.0.0.1:1071"
Environment="ASENA_GATEWAY_URL=http://127.0.0.1:8080"

# Ana çalıştırma komutu - 7/24 servis modu
ExecStart=/root/ASENA_CORE/target/release/asena serve --scan-interval 300 --scout --forge --self-healing

# Yeniden başlatma ayarları
Restart=always
RestartSec=15
StartLimitIntervalSec=60
StartLimitBurst=5

# Kaynak limitleri
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=4G
CPUQuota=200%

# Güvenlik ayarları
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/root/ASENA_CORE/data

# Log ayarları
StandardOutput=journal
StandardError=journal
SyslogIdentifier=asena-main

# Timeout ayarları
TimeoutStartSec=120
TimeoutStopSec=30

# Önbellek temizliği (çökme durumunda)
ExecStopPost=/bin/bash -c 'rm -f /root/ASENA_CORE/data/*.lock 2>/dev/null || true'

[Install]
WantedBy=multi-user.target
