# deploy/sentient-gateway.service
# ─── SENTIENT GATEWAY SERVICE ───
# SENTIENT OS - API Gateway
# HTTP REST API, WebSocket, Telegram Bot
#
# Kurulum:
#   sudo cp deploy/sentient-gateway.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable sentient-gateway
#   sudo systemctl start sentient-gateway

[Unit]
Description=SENTIENT Gateway - API & WebSocket Server
Documentation=https://github.com/nexsusagent-coder/SENTIENT_CORE
After=network.target network-online.target sentient-vgate.service
Wants=network-online.target
Requires=sentient-vgate.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/SENTIENT_CORE

# Ortam değişkenleri (.env dosyasından yüklenir)
EnvironmentFile=/root/SENTIENT_CORE/.env

# Çalıştırma komutu
ExecStart=/root/SENTIENT_CORE/target/release/sentient gateway --http-addr 0.0.0.0:8080

# Yeniden başlatma ayarları
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

# Kaynak limitleri
LimitNOFILE=65536
LimitNPROC=2048
MemoryMax=2G
CPUQuota=150%

# Güvenlik ayarları
NoNewPrivileges=true
PrivateTmp=true

# Log ayarları
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentient-gateway

# Timeout ayarları
TimeoutStartSec=60
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
